# deployment/isoproxy.service
[Unit]
Description=Isoproxy - Minimal Anthropic-compatible proxy for Claude Code
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=isoproxy
Group=isoproxy
WorkingDirectory=/opt/isoproxy

# Load environment variables from secure config file
# Create with: sudo -u isoproxy nano /etc/isoproxy/config.env
# Secure with: sudo chmod 600 /etc/isoproxy/config.env
EnvironmentFile=/etc/isoproxy/config.env

# Start the proxy server
ExecStart=/opt/isoproxy/.venv/bin/uvicorn isoproxy.main:app \
  --host 127.0.0.1 \
  --port 9000 \
  --workers 1

# Restart configuration
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Network restrictions
# Note: This allows IPv4/IPv6 but does NOT restrict destinations.
# Application-level endpoint allowlisting provides connection control.
RestrictAddressFamilies=AF_INET AF_INET6

# System call filtering
# @network-io is required for Python async runtimes (uvicorn, asyncio)
SystemCallFilter=@system-service @network-io
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=1024
CPUQuota=50%
MemoryMax=512M

[Install]
WantedBy=multi-user.target
