# deployment/isoproxy-socket.service
[Unit]
Description=Isoproxy - Minimal Anthropic-compatible proxy for Claude Code (Unix Socket)
After=network.target

[Service]
Type=simple
User=isoproxy
Group=isoproxy
WorkingDirectory=/opt/isoproxy
Environment=PYTHONPATH=/opt/isoproxy
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/isoproxy/config.env
RuntimeDirectory=isoproxy

# Start the proxy server with Unix socket
ExecStart=/opt/isoproxy/.venv/bin/uvicorn --app-dir /opt/isoproxy isoproxy.main:app --uds /run/isoproxy/isoproxy.sock --workers 1

# Restart configuration
Restart=on-failure
RestartSec=3

# Basic security hardening
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target